# PC_CONTROL_CODE/VerySimpleCode/crontab
# Run every minute
* * * * /location/to/bash/simple.sh